# salt/compunaut_consul/alerts/build/consul-alerts.service
[Unit]
Description=A highly available daemon for sending notifications and reminders based on Consul health checks.
Documentation=https://github.com/AcalephStorage/consul-alerts
After=network-online.target

[Service]
User=consul_alerts
Group=consul_alerts
LimitNOFILE=65536
ExecStart=/opt/consul_alerts/bin/consul-alerts -config /etc/consul-alerts/consul-alerts.conf
KillMode=control-group
Restart=on-failure

[Install]
WantedBy=multi-user.target
